# packaging/DEBIAN/prerm
#!/bin/bash
set -e

# Stop the service before removing files
if systemctl is-active --quiet fleet-agent; then
    systemctl stop fleet-agent
    echo "fleet-agent stopped"
fi
systemctl disable fleet-agent 2>/dev/null || true
