# docs/fleet-agent.service
[Unit]
Description=Fleet Manager Agent
Documentation=https://github.com/Lightbridge-KS/FleetManager
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/opt/fleet-agent/Agent.Worker
WorkingDirectory=/var/lib/fleet-agent
EnvironmentFile=/etc/fleet-agent/agent.env

User=fleet-agent
Group=fleet-agent

Restart=always
RestartSec=10

# Security hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/fleet-agent
ReadOnlyPaths=/opt/fleet-agent

[Install]
WantedBy=multi-user.target
